# arm_z.service
[Unit]
Description=Sub: 8410 - End effector target volocity

[Service]
ExecStart=/usr/bin/python /home/pi/RoverArm/z_controller.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
